# package/docker/anon/netrics/volume/etc/cron.d/cron-reboot
1 * * * * root /etc/init.d/nm-exp-active-netrics restart
