# package/deb/nm-exp-active-netrics/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/nm-exp-active-netrics
update-rc.d nm-exp-active-netrics defaults
/etc/init.d/nm-exp-active-netrics start
/usr/local/bin/netrics -C
